# deploy/beava.service
[Unit]
Description=Beava real-time feature server
Documentation=https://github.com/petrpan26/beava
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=beava
Group=beava
# --tcp-bind 127.0.0.1 is SECURITY-CRITICAL: the raw binary TCP protocol has no auth.
# Port 6400 MUST only be reachable from the local loopback interface. Defense in depth:
# ufw in provision.sh also denies 6400/tcp on the public interface.
ExecStart=/usr/local/bin/beava \
  --public-mode \
  --tcp-bind 127.0.0.1 \
  --tcp-port 6400 \
  --http-port 6401 \
  --snapshot-dir /var/lib/beava
Environment=RUST_LOG=info
# /etc/beava/admin.token.env contains BEAVA_ADMIN_TOKEN=<hex> (mode 600, owner beava:beava)
EnvironmentFile=/etc/beava/admin.token.env
Restart=always
RestartSec=5
# StateDirectory=beava manages /var/lib/beava with correct ownership for the snapshot dir.
StateDirectory=beava
StandardOutput=journal
StandardError=journal
LimitNOFILE=65536

# Sandboxing — cheap hardening for a public-facing service.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ReadWritePaths=/var/lib/beava

[Install]
WantedBy=multi-user.target
